# control/postinst
#!/bin/sh
systemctl enable buildkit.socket
systemctl enable buildkit.service
systemctl start buildkit.socket
