# control/prerm
#!/bin/sh
systemctl stop buildkit.socket
systemctl disable buildkit.service
systemctl disable buildkit.socket
